# etc/systemd/system/elasticsearch.service.d/override.conf
[Service]

# Sacrifice search for memory in dire times.
Restart = no
OOMScoreAdjust = 400

MemoryMax = 2G
MemoryHigh = 1.5G

LimitNOFILE = 65535:524288
